# src/systemd/dotfile-vcs-check.timer
[Unit]
Description=Check for changes in dotfiles

[Timer]
OnBootSec=10s
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
