# src/systemd/dotfile-vcs-check.service
[Unit]
Descritpion=Check for changes in dotfiles

[Service]
Type=oneshot
RemainAfterExit=yes
SuccessExitStatus=1
ExecStart=dotfile-vcs-check
